# raspi-pilot.service
[Unit]
Description=Raspberry Pi Rover

[Service]
Type=simple
ExecStart=raspi-pilot
Restart=always

[Install]
WantedBy=multi-user.target